# src/scripts/bay-systemd/systemd/cocalc-bay-cloudflared.service
[Unit]
Description=CoCalc Bay Cloudflare Tunnel
After=network-online.target cocalc-bay-frontdoor.service
Wants=cocalc-bay-frontdoor.service

[Service]
Type=exec
User=cocalc-bay
WorkingDirectory=/opt/cocalc/bay/current
EnvironmentFile=/etc/cocalc/bay.env
EnvironmentFile=-/etc/cocalc/bay-workers.env
EnvironmentFile=-/etc/cocalc/bay-overlay.env
EnvironmentFile=-/etc/cocalc/bay-topology.env
EnvironmentFile=-/etc/cocalc/bay-secrets.env
Environment=COCALC_REQUIRE_SITE_MASTER_KEY=1
Environment=COCALC_BAY_CLOUDFLARED_SYSTEMD=1
LoadCredential=site-master-key:/etc/cocalc/site-master-key
ExecStartPre=/opt/cocalc/bay/current/bin/bay-preflight cloudflared
ExecStart=/opt/cocalc/bay/current/bin/bay-cloudflared
Restart=on-failure
RestartSec=5
TimeoutStartSec=120
TimeoutStopSec=20
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
